# clns-acuity-transform/acuity-transform.cron
0 22 * * * cd ~/acuity-transform && ./acuity-transform.sh
